# Ghidra/jdk.jpackage/jdk/jpackage/internal/resources/unit-template.service
[Unit]
Description=SERVICE_DESCRIPTION

[Service]
ExecStart=APPLICATION_LAUNCHER
Restart=on-failure

[Install]
WantedBy=multi-user.target
